# config/includes.chroot/etc/systemd/system/neuraldrive-caddy.service
[Unit]
Description=NeuralDrive Caddy Reverse Proxy
After=network-online.target neuraldrive-certs.service
Requires=neuraldrive-certs.service
Wants=network-online.target

[Service]
User=neuraldrive-caddy
Group=neuraldrive-caddy
EnvironmentFile=/etc/neuraldrive/caddy.env
Environment=XDG_DATA_HOME=/var/lib/neuraldrive/caddy
Environment=XDG_CONFIG_HOME=/var/lib/neuraldrive/caddy/config
ExecStart=/usr/local/bin/caddy run --config /etc/neuraldrive/Caddyfile --adapter caddyfile
ExecReload=/usr/local/bin/caddy reload --config /etc/neuraldrive/Caddyfile --adapter caddyfile
Restart=always
RestartSec=5

AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
ReadOnlyPaths=/usr/local/bin/caddy
ReadWritePaths=/var/log/neuraldrive /var/lib/neuraldrive/caddy

[Install]
WantedBy=multi-user.target
