# config/includes.chroot/etc/systemd/system/neuraldrive-gpu-monitor.service
[Unit]
Description=NeuralDrive GPU Monitoring Service
After=neuraldrive-gpu-detect.service

[Service]
User=neuraldrive-monitor
Group=neuraldrive-monitor
SupplementaryGroups=video render
WorkingDirectory=/usr/lib/neuraldrive/gpu-monitor/gpu-hot
ExecStart=/usr/lib/neuraldrive/gpu-monitor/venv/bin/python app.py --port 1312 --host 127.0.0.1
Restart=always
RestartSec=5

ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=no
DeviceAllow=/dev/nvidia* r
DeviceAllow=/dev/dri/* r
ReadOnlyPaths=/
ReadWritePaths=/var/log/neuraldrive

[Install]
WantedBy=multi-user.target
