# config/includes.chroot/etc/systemd/system/neuraldrive-ollama.service
[Unit]
Description=NeuralDrive Ollama Service
After=network.target neuraldrive-gpu-detect.service
Requires=neuraldrive-gpu-detect.service

[Service]
EnvironmentFile=/etc/neuraldrive/ollama.conf
ExecStartPre=/usr/bin/mkdir -p /var/lib/neuraldrive/models
ExecStart=/usr/local/bin/ollama serve
User=neuraldrive-ollama
Group=neuraldrive-ollama
SupplementaryGroups=video render
Restart=always
RestartSec=5
KillMode=process

LimitNOFILE=65535
MemoryHigh=90%
MemoryMax=95%

ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=no
ProtectKernelTunables=yes
ProtectControlGroups=yes
DeviceAllow=/dev/nvidia* rw
DeviceAllow=/dev/dri/* rw
ReadOnlyPaths=/
ReadWritePaths=/var/lib/neuraldrive/models /var/log/neuraldrive /etc/neuraldrive /run/neuraldrive

[Install]
WantedBy=multi-user.target
