# config/includes.chroot/etc/systemd/system/neuraldrive-setup.service
[Unit]
Description=NeuralDrive Initial Setup
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/lib/neuraldrive/first-boot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
