# config/includes.chroot/etc/systemd/system/nftables.service.d/neuraldrive.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/nft -f /etc/neuraldrive/nftables.conf
ExecReload=
ExecReload=/usr/sbin/nft -f /etc/neuraldrive/nftables.conf
